# cookbooks/j2q/templates/default/j2q.service
[Unit]
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=<%= @path %>/current/j2q \
	-environment <%= node.chef_environment %> \
	-cluster <%= node.cluster_name %> \
	-kafka <%= kafka_connect(node.cluster_name) %> \
	-port 1040
User=root
Group=root
Restart=always
WorkingDirectory=<%= @path %>/current
SyslogIdentifier=%p

[Install]
WantedBy=multi-user.target
