# cookbooks/xinetd/files/default/xinetd.service
[Unit]
Description=xinetd
After=network.target

[Service]
ExecStart=/usr/sbin/xinetd -dontfork
Restart=always

[Install]
WantedBy=multi-user.target
