# cookbooks/linux/files/default/irqd.service
[Unit]
Description=CPU Interrupt Balancer

[Service]
ExecStart=/usr/sbin/irqd -fd

[Install]
WantedBy=multi-user.target
