# cookbooks/glassfish/files/default/glassfish.initd
#!/sbin/runscript

GF_USER="glassfish"
GF_HOME="/usr/local/glassfish"
ASADMIN="${GF_HOME}/bin/asadmin"
SU="su --login ${GF_USER}"

depend() {
	need net
}

start() {
	ebegin "Starting ${SVCNAME}"
	${SU} -c "${ASADMIN} start-domain"
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	${SU} -c "${ASADMIN} stop-domain"
	eend $?
}
