# cookbooks/camus/templates/default/camus.service
[Unit]
Description=Camus (Kafka -> Hadoop importer)

[Service]
ExecStart=/var/app/camus/current/camus
User=camus
Group=camus
SyslogIdentifier=camus

[Install]
WantedBy=default.target
