# cookbooks/monit/templates/default/monit.initd
#!/sbin/runscript

name="monit for user <%= @user[:name] %>"

pidfile=<%= @pidfile %>
command=/usr/bin/monit
command_args="-p ${pidfile}"

command_user=<%= @user[:name] %>

depend() {
	use net
}
