# cookbooks/zenops/files/default/pkgsync.service
[Unit]
Description=Package Sync
After=network.target

[Service]
ExecStart=/usr/sbin/pkgsync
SyslogIdentifier=pkgsync
