# cookbooks/zenops/files/default/pkgsync.timer
[Unit]
Description=Package Sync Timer

[Timer]
Unit=pkgsync.service
OnCalendar=*:25

[Install]
WantedBy=multi-user.target
