# cookbooks/nepal/files/default/nepald.service
[Unit]
Description=nepal Management Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/nepald
Restart=always
OOMScoreAdjust=-900

[Install]
WantedBy=multi-user.target
