# cookbooks/chef/files/default/chef-client.timer
[Unit]
Description=Chef Client Timer

[Timer]
Unit=chef-client.service
OnBootSec=3600
OnUnitInactiveSec=3600

[Install]
WantedBy=multi-user.target
